# pushgit.timer
[Unit]
Description=Run pushgit.service daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=pushgit.service

[Install]
WantedBy=timers.target