# files/pushgit.timer
[Unit]
Description=Run pushgit.service

[Timer]
OnCalendar=*:0/10
Persistent=true
Unit=pushgit.service

[Install]
WantedBy=timers.target
